# supplementary/skaffari.service.in
[Unit]
Description=Email account management web interface
Documentation=man:skaffari(8) man:skaffari.ini(5)
After=network.target
After=mysql.service cyrus-imapd.service postfix.service
Before=httpd.service nginx.service lighttpd.service

[Service]
Type=notify
Environment="QT_LOGGING_CONF=@SKAFFARI_INI_FILE@"
PIDFile=@SKAFFARI_PID_FILE@
ExecStartPre=@SKAFFARI_CMD_PATH@ -t -i @SKAFFARI_INI_FILE@
ExecStart=@CUTELYST_WSGI_EXECUTABLE@ -M -a @SKAFFARI_APP_FILE@ --ini @SKAFFARI_INI_FILE@ --pidfile @SKAFFARI_PID_FILE@ --fastcgi-socket @SKAFFARI_SOCKET_FILE@ --socket-access ug
ExecStop=@CUTELYST_WSGI_EXECUTABLE@ --stop @SKAFFARI_PID_FILE@
NotifyAccess=all
RuntimeDirectory=skaffari
PrivateTmp=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
User=@SKAFFARI_USER@
Group=@SKAFFARI_GROUP@

[Install]
WantedBy=multi-user.target
